# rc.local
#!/bin/sh
# ip link set can0 up type can bitrate 250000
ip link set can0 up type can bitrate 250000
svc -u /service/dbus-ubms.can0
