# rc.local
#!/bin/sh
ip link set can0 up type can bitrate 250000
svc -o /service/dbus-ubms
